# deploy/trading-bot.service
[Unit]
Description=Trading Bot Manager - Sistema 24/7
After=network.target
Wants=network.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i/trading-bot
Environment="PATH=/home/%i/trading-bot/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/%i/trading-bot"
EnvironmentFile=/home/%i/trading-bot/.env
ExecStart=/home/%i/trading-bot/venv/bin/python /home/%i/trading-bot/bot.py

Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

TimeoutStopSec=30
KillMode=mixed
FinalKillSignal=SIGKILL

StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=trading-bot

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/%i/trading-bot/data /home/%i/trading-bot/logs /home/%i/trading-bot/dashboard/output

[Install]
WantedBy=multi-user.target